# templates/node_exporter.service.epp
<%- |
  String $install_path,
  String $version,
  String $flavor,
  | -%>

[Unit]
Description=Node Exporter

[Service]
User=<%= $usermame %>
EnvironmentFile=<%= $config_file %>
ExecStart=<%= $install_path %>/node_exporter-<%= $version %>.<%= $flavor %>/node_exporter $ARGS

[Install]
WantedBy=multi-user.target
